# DEBIAN/postinst
# Installare Suricata se non è già installato
if ! dpkg -s suricata >/dev/null 2>&1; then
    apt update
    apt install -y suricata
fi

# Copia del file suricata.yml personalizzato
install -m 644 etc/suricata/suricata.yml /etc/suricata/suricata.yml

# Copia del file suricata.service personalizzato
install -m 644 lib/systemd/system/suricata.service /lib/systemd/system/suricata.service

# Ricaricare il daemon di systemd per rilevare le modifiche al servizio
systemctl daemon-reload

# Avviare il servizio Suricata
systemctl start suricata

# Abilitare Suricata all'avvio del sistema
systemctl enable suricata

exit 0
